# crontab
30 20 * * * /usr/bin/python /home/ubuntu/pyquant/update_historical_prices.py
