# crontab
MAILTO=norbinn@gmail.com

52 18 * * * . $HOME/.profile; /usr/bin/python /home/ubuntu/pyquant/update_historical_prices.py
