# crontab
30 13 * * * /usr/bin/python /home/ubuntu/pyquant/update_historical_prices.py
